# vbb-rest/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio

bashio::log.info "Waiting for vbb-rest to start before starting nginx"
bashio::net.wait_for 3000 localhost 300
bashio::log.info "Starting nginx"
exec nginx
